# TwinklyWall/ddp_bridge.service
[Unit]
Description=TwinklyWall DDP to Pixel Overlay Bridge
After=network.target

[Service]
Type=simple
User=fpp
Group=fpp
WorkingDirectory=/home/fpp/TwinklyWall_Project/TwinklyWall
Environment="PATH=/home/fpp/TwinklyWall_Project/TwinklyWall/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="FPP_MODEL_NAME=Light Wall"
ExecStart=/home/fpp/TwinklyWall_Project/TwinklyWall/.venv/bin/python ddp_bridge.py --port 4049 --width 90 --height 50 --model "Light_Wall"
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
